# systemd/gravedecay-term.service.tmpl
[Unit]
Description=gravedecay web terminal (ttyd -> tmux agents socket)
After=network-online.target
Wants=network-online.target
# Page the human if this unit lands in failed state (silent no-op unless
# notifications are configured -- see docs/NOTIFICATIONS.md).
OnFailure=gravedecay-notify@%n.service

[Service]
Type=simple
User=@USER@
Group=@USER@
# -i lo: loopback only (reached via tailscale serve /term).  -W: writable.
# --url-arg lets the launcher tiles pick the session (?arg=claude etc).
# No --base-path: tailscale serve strips the /term mount prefix.
# -I: custom frontend with working clipboard (OSC 52 + copy/paste bar) —
#     built by raise.sh from web/term/, see docs/TERMINAL.md.
ExecStart=@TTYD@ -p @TERM_PORT@ -i lo -W --url-arg \
  -I @GRAVE_ROOT@/web/term/index.html \
  -t fontSize=14 -t 'theme={"background":"#0d0d0d"}' -t titleFixed=gravedecay \
  @GRAVE_ROOT@/scripts/webterm
Restart=always
RestartSec=3
SyslogIdentifier=gravedecay-term
# Agent CLIs (claude/codex in ~/.local/bin) must be reachable from sessions.
# @TOOLPATH@ prepends the durable toolchain when it lives under $HOME
# (Homebrew node/bins, rootless docker) — empty on package-manager hosts.
Environment=PATH=@TOOLPATH@@HOME@/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=@HOME@
Environment=GRAVE_ROOT=@GRAVE_ROOT@
Environment=DOCKER_HOST=@DOCKER_HOST@

# NO NoNewPrivileges: it blocks sudo — the dashboard runs mode flips via
# sudo -n grave, and web-terminal shells legitimately use sudo.
# NO PrivateTmp: tmux sockets live in /tmp/tmux-<uid>/ — a private /tmp
# would hide agent sessions (dashboard) / strand them in a namespace (ttyd)

[Install]
WantedBy=multi-user.target
