# systemd/gravedecay.service.tmpl
[Unit]
Description=gravedecay - appliance dashboard & entry point
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@USER@
Group=@USER@
ExecStart=/usr/bin/python3 @GRAVE_ROOT@/scripts/gravedecay.py
Restart=always
RestartSec=3
SyslogIdentifier=gravedecay
Environment=GRAVE_ROOT=@GRAVE_ROOT@
Environment=GRAVEDECAY_PORT=@DASH_PORT@
# Comma-separated Tailscale logins allowed to press action buttons.
# Empty = tailnet viewers are read-only. Add yours in a drop-in or here.
Environment=GRAVEDECAY_ALLOWED_USERS=@ALLOWED_USERS@
# Launcher tiles ("label=url;label=url") — gravedecay is the entry point;
# every app mounted on this box gets a tile here.
Environment="GRAVEDECAY_APPS=⌨️ T3 Code=/;🖥️ Terminal=/term/?arg=shell;🤖 Claude=/term/?arg=claude;🧠 Codex=/term/?arg=codex"

# NO NoNewPrivileges: it blocks sudo — the dashboard runs mode flips via
# sudo -n grave, and web-terminal shells legitimately use sudo.
# NO PrivateTmp: tmux sockets live in /tmp/tmux-<uid>/ — a private /tmp
# would hide agent sessions (dashboard) / strand them in a namespace (ttyd)

[Install]
WantedBy=multi-user.target
